# systemd/prusa-slack-bot.service
# Sample systemd unit for Path B (native install on Raspberry Pi OS / Debian).
#
# Install:
#   sudo cp systemd/prusa-slack-bot.service /etc/systemd/system/
#   sudo cp .env.example /etc/prusa-slack-bot.env  # then edit it
#   sudo chmod 600 /etc/prusa-slack-bot.env
#   sudo useradd --system --home /var/lib/prusa-slack-bot --create-home prusabot
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now prusa-slack-bot.service
#
# This unit assumes you've installed the package as user `prusabot` with
# something like `pipx install prusa-slack-bot` or via a venv at
# /var/lib/prusa-slack-bot/.venv. Adjust ExecStart to taste.

[Unit]
Description=Prusa printer Slack bot
Documentation=https://github.com/your-org/prusa-slack-bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=prusabot
Group=prusabot
WorkingDirectory=/var/lib/prusa-slack-bot
EnvironmentFile=/etc/prusa-slack-bot.env
Environment=DB_PATH=/var/lib/prusa-slack-bot/prusa-slack-bot.sqlite3
ExecStart=/var/lib/prusa-slack-bot/.venv/bin/python -m prusa_slack_bot
Restart=on-failure
RestartSec=10s
KillSignal=SIGTERM
TimeoutStopSec=15s

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/prusa-slack-bot
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
CapabilityBoundingSet=
AmbientCapabilities=

[Install]
WantedBy=multi-user.target
